# etc/systemd/system/upmpdcli.service.d/override.conf
[Service]
EnvironmentFile=-/etc/default/upmpdcli

