# etc/systemd/system/squeezelite.service.d/override.conf
[Service]
LimitNOFILE=32768
ExecStartPost=-/usr/local/bin/optimize_player squeezelite 4

